# nginx-php-laravel-one/master.ini

[supervisord]
logfile=/tmp/logs/supervisord.log
logfile_maxbytes=50MB        
logfile_backups=10           
loglevel=info                
pidfile=/tmp/supervisord.pid 
nodaemon=true   
user=root            
minfds=1024                 
minprocs=200                

[program:crond]
command=crond -f -L 15
user=root
autostart=true

[program:php-fpm]
command = /usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
stderr_logfile=/tmp/logs/supervisor-php.log
stdout_logfile=/tmp/logs/supervisor-php.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/tmp/logs/supervisor-nginx.log
stdout_logfile=/tmp/logs/supervisor-nginx.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[include]
files=/etc/supervisor.d/*.ini

[unix_http_server]
file=/tmp/supervisor.sock
